# HOME/.config/autostart/compton.desktop
[Desktop Entry]
Name=Compton
Exec=compton
Terminal=false
Type=Application
