# HOME/.config/autostart/signal.desktop
[Desktop Entry]
Name=Signal
Exec=/home/fred/.bin/signal.sh
Terminal=false
Type=Application
