# HOME/.config/autostart/slack.desktop
[Desktop Entry]
Name=Slack
Exec=slack
Terminal=false
Type=Application
