# HOME/.config/autostart/slack.desktop
[Desktop Entry]
Name=Tresorit
Exec=slack
Terminal=false
Type=Application
