# HOME/.config/autostart/trayer.desktop
[Desktop Entry]
Name=Trayer
Exec=/home/fred/bin/run_trayer.sh
Terminal=false
Type=Application
