# HOME/.config/autostart/tresorit.desktop
[Desktop Entry]
Name=Tresorit
Exec=/home/fred/.bin/tresorit.sh
Terminal=false
Type=Application
